# ARM/systemd/lapee-hyperbeam.service
[Unit]
Description=LapEE ARM HyperBEAM node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/lapee-arm/hyperbeam
Environment=HB_CONFIG=/etc/lapee-arm/lapee-arm.json
Environment=HB_MODE=debug
Environment=LAPEE_TPM_ALLOW_NO_NIF=1
Environment=LAPEE_TPM_TCTI=device:/dev/tpm0
ExecStart=/opt/lapee-arm/hyperbeam/bin/hb foreground
Restart=no
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
